# conteur/deploy/conteur-bucephale-sync.service
[Unit]
Description=Sync live Bucéphale chapters for Conteur

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/conteur-sync-bucephale
